# gunicorn.supervisor.conf
[program:gunicorn]

directory=/app/djangoapp
command=/usr/local/bin/gunicorn -c /app/gunicorn_conf.py SenzWeb.wsgi:application
stdout_logfile=/app/logs/gunicorn-stdout.log
stdout_logfile_maxbytes=1MB
stderr_logfile=/app/logs/gunicorn-stderr.log
stderr_logfile_maxbytes=1MB
stopsignal=KILL
killasgroup=true
autostart=true
autorestart=true

